# msi-control/contents/scripts/99-msi-ec.rules
# MSI-EC platform device attributes
# Grant rw-rw---- (0660) to 'msi-ec' group for CONTROL nodes
# NOTE: ATTR{x}="value" writes string to sysfs, it does NOT chmod.
#       For sysfs permission changes, RUN+= with chmod/chgrp is required.

SUBSYSTEM=="platform", DRIVERS=="msi-ec", \
  RUN+="/bin/sh -c '\
    cd /sys/devices/platform/msi-ec && \
    for f in shift_mode fan_mode cooler_boost webcam webcam_block super_battery fn_key win_key usb_power; do \
      [ -e \"$f\" ] && chgrp msi-ec \"$f\" && chmod 0660 \"$f\"; \
    done; \
    for f in cpu/basic_fan_speed gpu/basic_fan_speed; do \
      [ -e \"$f\" ] && chgrp msi-ec \"$f\" && chmod 0660 \"$f\"; \
    done \
  '"

# Battery threshold control attributes
SUBSYSTEM=="power_supply", KERNEL=="BAT*", \
  RUN+="/bin/sh -c '\
    cd /sys%p && \
    for f in charge_control_end_threshold charge_control_start_threshold; do \
      [ -e \"$f\" ] && chgrp msi-ec \"$f\" && chmod 0660 \"$f\"; \
    done \
  '"

# Keyboard backlight brightness
SUBSYSTEM=="leds", KERNEL=="msiacpi::kbd_backlight", \
  RUN+="/bin/sh -c '\
    [ -e /sys%p/brightness ] && chgrp msi-ec /sys%p/brightness && chmod 0660 /sys%p/brightness \
  '"
